# Centos7-Pmaker/files/laravel-worker-workflow.ini
[program:laravel-worker-workflow]
user=apache
directory=/opt/processmaker/
command=/opt/processmaker/processmaker artisan queue:work --workspace=workflow
stdout_logfile=/opt/processmaker/worker-workflow.log
autostart=true
autorestart=true